# files/etc/uci-defaults/99-custom
#!/bin/sh
# First-boot customization

# Set DNS
uci set network.lan.dns='127.0.0.1'
uci commit network

exit 0
